# etc/cron/crontab
MAILTO=150283t@gmail.com

05 * * * * /usr/bin/docker system prune -f >/dev/null
