# etc/cron/crontab

05 * * * * /usr/bin/docker system prune
