# etc/cron/crontab
MAILTO=150283t@gmail.com

30 18 * * * /usr/bin/docker system prune -f >/dev/null
30 18 * * * mv /Users/{userName}/Downloads/* /tmp
